# ops/systemd/clearprice-worker.service
[Unit]
Description=ClearPrice Celery Worker
After=network.target redis-server.service postgresql.service
Wants=redis-server.service postgresql.service

[Service]
Type=exec
User=clearprice
Group=clearprice
WorkingDirectory=/opt/clearprice/backend
Environment="PATH=/opt/clearprice/backend/venv/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/opt/clearprice/backend"
EnvironmentFile=/opt/clearprice/backend/.env
ExecStart=/opt/clearprice/backend/venv/bin/celery -A app.tasks.celery_app worker --loglevel=info --concurrency=4
Restart=always
RestartSec=10
StandardOutput=append:/opt/clearprice/logs/worker.log
StandardError=append:/opt/clearprice/logs/worker.log

[Install]
WantedBy=multi-user.target
