# ops/systemd/clearprice-api.service
[Unit]
Description=ClearPrice FastAPI Backend
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
Type=exec
User=clearprice
Group=clearprice
WorkingDirectory=/opt/clearprice/backend
Environment="PATH=/opt/clearprice/backend/venv/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/opt/clearprice/backend"
EnvironmentFile=/opt/clearprice/backend/.env
ExecStart=/opt/clearprice/backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 4
Restart=always
RestartSec=5
StandardOutput=append:/opt/clearprice/logs/api.log
StandardError=append:/opt/clearprice/logs/api.log

[Install]
WantedBy=multi-user.target
