# ops/systemd/clearprice-flower.service
[Unit]
Description=ClearPrice Flower Celery Monitor
After=network.target redis-server.service
Wants=redis-server.service

[Service]
Type=exec
User=clearprice
Group=clearprice
WorkingDirectory=/opt/clearprice/backend
Environment="PATH=/opt/clearprice/backend/venv/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/opt/clearprice/backend"
EnvironmentFile=/opt/clearprice/backend/.env
ExecStart=/opt/clearprice/backend/venv/bin/celery -A app.tasks.celery_app flower --address=127.0.0.1 --port=5555 --basic_auth=admin:REPLACE_WITH_FLOWER_PASSWORD
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
